# bin/tinymonitor/sample.cron
*/5 * * * * /usr/local/tinymonitor/monitor
